# Daemon/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.chargecontrol.daemon</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/ChargeControlDaemon</string>
	<key>AssociatedBundleIdentifiers</key>
	<string>com.chargecontrol.app</string>
	<key>MachServices</key>
	<dict>
		<key>com.chargecontrol.daemon</key>
		<true/>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
